# roles/node_exporter/templates/node-exporter.service.j2
[Unit]
Description=Prometheus Node Exporter
After=network.target

[Service]
Type=simple
User=node_exporter
Group=node_exporter
ExecStart=/usr/local/bin/node_exporter \
    --web.listen-address={{ node_exporter_host }}:{{ node_exporter_port }}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
